# packaging/debian/postinst
#!/bin/sh
set -e

if [ "$1" = "configure" ]; then
    systemctl daemon-reload || true
    systemctl enable --now linuxnetwatch-collector.service || true
fi

exit 0
