# packaging/debian/prerm
#!/bin/sh
set -e

if [ "$1" = "remove" ] || [ "$1" = "purge" ]; then
    systemctl disable --now linuxnetwatch-collector.service || true
fi

exit 0
